# iperf3/server/iperf3.service
[Unit]
Description=iperf3 server
After=syslog.target network.target auditd.service

[Service]
ExecStart=/usr/bin/iperf3 -s --port 5201

[Install]
WantedBy=multi-user.target